# systemd/user/sol-notes-snapshot.service
[Unit]
Description=Run sol-notes-snapshot

[Service]
ExecStart=zsh -lc sol-notes-snapshot
SyslogIdentifier=sol-notes-snapshot
